# deploy/systemd/cat-monitor.service
[Unit]
Description=cat-monitor backend
After=network-online.target mosquitto.service zigbee2mqtt.service
Wants=network-online.target mosquitto.service

[Service]
Type=simple
WorkingDirectory=/opt/cat-monitor
EnvironmentFile=/opt/cat-monitor/.env
ExecStart=/opt/cat-monitor/backend/target/release/cat-monitor-rust-backend
Restart=on-failure
RestartSec=5
User=catmonitor
Group=catmonitor

[Install]
WantedBy=multi-user.target
