# deploy/systemd/cloudflared-cat-monitor.service
[Unit]
Description=cat-monitor cloudflared tunnel
After=network-online.target cat-monitor.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cat-monitor
EnvironmentFile=/opt/cat-monitor/.env
ExecStart=/bin/sh -lc 'test -n "${CLOUDFLARE_TUNNEL_TOKEN:-}" && exec /usr/bin/env cloudflared tunnel --no-autoupdate --protocol "${CLOUDFLARED_PROTOCOL:-http2}" run'
Restart=on-failure
RestartSec=5
User=catmonitor
Group=catmonitor

[Install]
WantedBy=multi-user.target
